# com.goldbot.runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.goldbot.runner</string>

    <!-- caffeinate -i mencegah Mac tidur; python3 yang punya pandas/sklearn -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>/usr/bin/python3</string>
        <string>live_runner.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yopiangga/Documents/Riset/64 Bot Forex/4</string>

    <!-- Mulai otomatis saat login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart otomatis kalau crash/keluar -->
    <key>KeepAlive</key>
    <true/>

    <!-- Log keluaran -->
    <key>StandardOutPath</key>
    <string>/Users/yopiangga/Documents/Riset/64 Bot Forex/4/runner.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yopiangga/Documents/Riset/64 Bot Forex/4/runner.log</string>
</dict>
</plist>
